# templates/kresd.service.j2
############################################################################
########## Managed by ansible (role: dns-resolver), do not edit! ###########
############################################################################

[Unit]
Description=Knot DNS Resolver daemon
Documentation=man:kresd(8)
## This is NOT a socket-activated service:
## RefuseManualStart=true

[Service]
Type=simple
WorkingDirectory=/run/knot-resolver/cache
EnvironmentFile=/etc/default/kresd
ExecStart=/usr/sbin/kresd $DAEMON_ARGS
# We need to do this in the unit file, because we need to drop privileges after
# opening the ports :(
# User=knot-resolver
Restart=on-failure

[Install]
# WantedBy=sockets.target
WantedBy=multi-user.target
